# voice-gateway.service
[Unit]
Description=AI Voice Gateway Service
After=network.target

[Service]
Type=simple
User=ztmdsbt
WorkingDirectory=/home/ztmdsbt/workspaces/audio-proxy/wsl2
Environment="PATH=/home/ztmdsbt/workspaces/audio-proxy/venv/bin:/usr/bin:/bin"
ExecStart=/home/ztmdsbt/workspaces/audio-proxy/venv/bin/python3 /home/ztmdsbt/workspaces/audio-proxy/wsl2/agent-gateway.py
Restart=always
RestartSec=5
StandardOutput=append:/home/ztmdsbt/workspaces/audio-proxy/logs/voice-gateway.log
StandardError=append:/home/ztmdsbt/workspaces/audio-proxy/logs/voice-gateway-error.log

[Install]
WantedBy=multi-user.target
